# whats-spinning.service
[Unit]
Description=Whats Spinning - Automatic Music Recognition
After=network.target sound.target systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
Type=simple
User=pi
Group=audio
SupplementaryGroups=audio
WorkingDirectory=/home/pi/dev/whats-spinning
ExecStartPre=/bin/sleep 5
ExecStart=/home/pi/dev/whats-spinning/venv/bin/python -m src.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="AUDIODEV=hw:1,0"
Environment="PULSE_SERVER="

[Install]
WantedBy=multi-user.target
